# filesystems/BTRFS/btrfs-scrub.service
[Unit]
Description=Run btrfs scrub on filesystem

[Service]
Type=simple
ExecStart=/opt/btrfsScripts/btrfs_scrub.sh
