# filesystems/BTRFS/btrfs-scrub.timer
[Unit]
Description=Monthly btrfs scrub

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
